# salt/roots/salt/pydashery/supervisor-pydashery.conf
[program:pydashery]
command=bash /home/pi/pydashery/start_backend.sh
numprocs=1
user=pi
stopsignal=INT
autostart=true
autorestart=true
startretries=999999999
